# fabric_data/config_files_templates/docker/supervisord.conf
[supervisord]
nodaemon=true 

[program:mysql]
process_name=mysqld
command = /usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
redirect_stderr = true
stdout_logfile = /var/log/mysqld.log
stdout_events_enabled = true
startsecs = 30
autorestart = true
user_name=mysql

[program:php]
command = /usr/sbin/php-fpm
user = root
autostart = true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:sshd]
autostart=true
autorestart=true
command=/usr/sbin/sshd -D

[program:crond]
command=/bin/bash -c "source /etc/sysconfig/crond && exec /usr/sbin/crond -n $CRONDARGS"
process_name=%(program_name)s
numprocs=1
directory=/tmp
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=root
